# install/linux/bop-merge-gate.service
[Unit]
Description=bop merge-gate (one-shot)

[Service]
Type=oneshot
ExecStart=/usr/local/bin/bop merge-gate --once
