# launchd/com.yourorg.jobcard.roadmap-inbox.plist.sample
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>sh.bop.roadmap-inbox</string>

  <key>QueueDirectories</key>
  <array>
    <string>__REPO_ROOT__/examples/roadmap-inbox/drop</string>
  </array>

  <key>RunAtLoad</key>
  <true/>

  <key>ProgramArguments</key>
  <array>
    <string>nu</string>
    <string>__REPO_ROOT__/scripts/ingest_roadmap_hotfolder.nu</string>
    <string>--inbox</string>
    <string>__REPO_ROOT__/examples/roadmap-inbox/drop</string>
    <string>--cards-dir</string>
    <string>__REPO_ROOT__/.cards</string>
  </array>

  <key>StandardOutPath</key>
  <string>/tmp/bop-roadmap-inbox.log</string>

  <key>StandardErrorPath</key>
  <string>/tmp/bop-roadmap-inbox.err</string>
</dict>
</plist>
