# launchd/com.yourorg.bop.dispatcher.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.yourorg.bop.dispatcher</string>
  
  <key>ProgramArguments</key>
  <array>
    <string>/usr/local/bin/bop</string>
    <string>dispatcher</string>
    <string>--adapter</string>
    <string>adapters/mock.nu</string>
    <string>--max-workers</string>
    <string>3</string>
    <string>--poll-ms</string>
    <string>500</string>
    <string>--max-retries</string>
    <string>3</string>
    <string>--reap-ms</string>
    <string>1000</string>
  </array>
  
  <key>WorkingDirectory</key>
  <string>REPLACE_WITH_REPO_ROOT</string>
  
  <key>EnvironmentVariables</key>
  <dict>
    <key>CARDS_DIR</key>
    <string>REPLACE_WITH_REPO_ROOT/.cards</string>
  </dict>
  
  <key>KeepAlive</key>
  <true/>

  <key>RunAtLoad</key>
  <true/>

  <!-- Wake dispatcher immediately when cards land in pending/ or drafts/ -->
  <key>WatchPaths</key>
  <array>
    <string>REPLACE_WITH_REPO_ROOT/.cards/pending</string>
    <string>REPLACE_WITH_REPO_ROOT/.cards/drafts</string>
  </array>
  
  <key>StandardOutPath</key>
  <string>/tmp/bop-dispatcher.log</string>
  
  <key>StandardErrorPath</key>
  <string>/tmp/bop-dispatcher.err</string>
  
  <key>HardResourceLimits</key>
  <dict>
    <key>NumberOfFiles</key>
    <integer>1024</integer>
  </dict>
  
  <key>SoftResourceLimits</key>
  <dict>
    <key>NumberOfFiles</key>
    <integer>512</integer>
  </dict>
</dict>
</plist>
